# overlay/usr/lib/systemd/system/atomtap-forward.service
[Unit]
Description=AtomTap Ethernet-to-WiFi traffic mirror
After=network.target atomtap-firstboot-setup.service
Wants=atomtap-firstboot-setup.service
ConditionPathExists=/var/lib/atomtap/firstboot.done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/atomtap-forward.sh start
ExecStop=/usr/local/sbin/atomtap-forward.sh stop
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
