# overlay/usr/lib/systemd/system/atomtap-firstboot-setup.service
[Unit]
Description=AtomTap first-boot setup gate
After=multi-user.target
Before=atomtap-forward.service
ConditionPathExists=/var/lib/atomtap/firstboot.done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/true

[Install]
WantedBy=multi-user.target
