# overlay/usr/lib/systemd/system/atomtap-forward.service
[Unit]
Description=AtomTap Ethernet-to-WiFi traffic mirror
After=network-online.target
Wants=network-online.target
ConditionPathExists=/var/lib/atomtap/firstboot.done

[Service]
Type=exec
ExecStart=/usr/local/sbin/atomtap-forward.sh start
ExecStop=/usr/local/sbin/atomtap-forward.sh stop
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
